# railgun-blockip.service
[Service]
ExecStart=/srv/railgun/blockip/blockip.py

[Install]
WantedBy=multi-user.target
